# gitlab-server/usr/share/gitlab-install/systemd-units/gitlab-backup.timer
[Unit]
Description=GitLab Backup Timer

[Install]
WantedBy=gitlab.target

[Timer]
OnCalendar=*-*-* 00:*:*
