# openstack/usr/lib/systemd/system/openstack-ironic-conductor.service
[Unit]
Description=OpenStack Bare Metal Provisioning Service (code-named Ironic) Conductor server
ConditionPathExists=/etc/ironic/ironic.conf
After=syslog.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ironic
ExecStart=/usr/bin/ironic-conductor --config-file /etc/ironic/ironic.conf

[Install]
WantedBy=multi-user.target

