# openstack/usr/lib/systemd/system/openstack-nova-setup.service
[Unit]
Description=Run nova-setup Ansible scripts
After=local-fs.target libvirtd.service openstack-keystone-setup.service postgres-server.service

[Service]
ExecStart=/usr/bin/ansible-playbook -v -M /usr/share/ansible/ansible-openstack-modules -i /usr/share/openstack/hosts /usr/share/openstack/nova.yml

[Install]
WantedBy=multi-user.target
